# raspi/dpkg/DEBIAN/postinst
#!/bin/sh

# Reload systemd manager configuration and enable service
systemctl daemon-reload
systemctl enable ugglan.service
