# drone/dpkg/DEBIAN/prerm
#!/bin/sh

# Stop and disable service
systemctl stop ugglan.service
systemctl disable ugglan.service
